# systemd/ctfpwn.target
[Unit]
Description=Offensive CTF framework